# supervisor/3-php.conf
[program:php-fpm7.0]
command = /usr/sbin/php-fpm7.0
autostart = true
autorestart=true
priority=2
